# truncate_all_tables.service
[Unit]
Description=Truncate all HCD and Presto tables for a fresh start
After=hcd.service


[Service]
Type=simple
Environment=PATH=/home/watsonx/affiliate-junction-demo/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/home/watsonx/affiliate-junction-demo/.venv
ExecStartPre=/bin/sleep 120
ExecStart=/home/watsonx/affiliate-junction-demo/.venv/bin/python /home/watsonx/affiliate-junction-demo/truncate_all_tables.py
Restart=on-failure
User=watsonx
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
